# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
/usr/bin/bled -r
bash /etc/custom_service/start_service.sh
ttyd -p 7681 /bin/bash --login &
(sleep 120 && /usr/bin/autotimesync.sh) &
(sleep 120 && screen -dmS antidc ping google.com) &
exit 0
